# systemd/piped-healthcheck.service
[Unit]
Description=Piped healthcheck (web + Postgres + HTTPS)
After=network.target piped-web.service

[Service]
Type=oneshot
User=ubuntu
ExecStart=/home/ubuntu/piped/scripts/healthcheck.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=piped-healthcheck
