# systemd/piped-optimize.service
[Unit]
Description=Piped hourly ROAS optimization tick
After=network.target piped-web.service

[Service]
Type=oneshot
User=ubuntu
EnvironmentFile=/home/ubuntu/piped/.env.local
ExecStart=/usr/bin/curl --silent --show-error --fail \
    --max-time 540 \
    -H "Authorization: Bearer ${CRON_SECRET}" \
    http://127.0.0.1:3001/api/cron/optimize
StandardOutput=journal
StandardError=journal
SyslogIdentifier=piped-optimize
